# rotated8-dot-net.conf
description 'uWSGI server for Rotated8.net'

start on runlevel [2345]
stop on runlevel [!2345]

exec uwsgi /var/www/rotated8.net/uwsgi.ini
