# rotated8-dot-net.conf
description 'uWSGI server for Rotated8.net'

start on runlevel [2345]
stop on runlevel [!2345]

setuid rotated8
setgid rotated8

exec uwsgi --ini /var/www/rotated8.net/uwsgi.ini
